# docker-images/mysql/files/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:mysql]
priority = 10
command = /usr/bin/mysqld --user=mysql

[program:telegraf]
priority = 10
command = /root/telegraf/telegraf --config /root/config/telegraf.conf