# docker-images/mysql/files/supervisord.conf
[unix_mysql_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:mysql]
priority = 10
command = /usr/bin/mysqld --user=mysql
stderr_logfile = /var/log/mysql/mysql-stderr.log
stdout_logfile = /var/log/mysql/mysql-stdout.log

[program:telegraf]
priority = 20
command = /root/telegraf/telegraf --config /root/config/telegraf.conf
stderr_logfile = /var/log/telegraf/telegraf-stderr.log
stdout_logfile = /var/log/telegraf/telegraf-stdout.log

[program:fluentd]
priority = 30
command = fluentd -c /root/config/fluentd.conf -p /fluentd/plugins